# com.pga.board.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pga.board.update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/thomaswilliams/Desktop/CLAUDE/pga-betting-board/weekly-update.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/thomaswilliams/Desktop/CLAUDE/pga-betting-board</string>

    <!-- Every Monday at 07:05. If the Mac is asleep, launchd runs it on next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/thomaswilliams/Desktop/CLAUDE/pga-betting-board/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/thomaswilliams/Desktop/CLAUDE/pga-betting-board/launchd.err.log</string>
</dict>
</plist>
